# etc/init/anc-server.conf
description "anc-server"
author "Behemot behemot@yandex-team.ru"

start on runlevel [3]
stop on runlevel [0156]

exec /usr/bin/anc-server.py
respawn